# etc/peer_flask.conf
[program:peer_flask]
command = gunicorn wsgi_gunicorn:app -b 127.0.0.1:9000 
directory = /home/finn/www/peer_flask
user = finn